# selfbot-gpt.service
[Unit]
Description=GPT4 Telegram Selfbot
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/gpt-selfbot
ExecStart=/usr/bin/python3 /path/to/gpt-selfbot/bot.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
